# projects/NextOS/packages/sysutils/powerstate/system.d/powerstate.service
[Unit]
Description=AC/Battery Performance Service
RequiresMountsFor=/storage
Before=nextos.target
StartLimitIntervalSec=0

[Service]
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/powerstate
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
