# projects/NextOS/packages/sysutils/systemd/system.d/debugconfig.service
[Unit]
Description=Setup debug config
DefaultDependencies=no
After=systemd-tmpfiles-setup.service
RequiresMountsFor=/storage
ConditionKernelCommandLine=|debugging
ConditionPathExists=|/storage/.cache/debug.nextos

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/debugconf/*.conf /run/nextos/debug'
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
