# projects/NextOS/packages/tools/inputplumber/sources/usr/lib/systemd/system/inputplumber-refresh.service
[Unit]
Description=Restart InputPlumber to pick up new HID devices

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart inputplumber.service
