# projects/NextOS/packages/sysutils/system-utils/system.d/headphones.service
[Unit]
Description=Headphone Sense
Before=nextos.target
After=pipewire.service
RequiresMountsFor=/storage

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/headphone_sense
Restart=on-failure

[Install]
WantedBy=multi-user.target
