# projects/NextOS/devices/Amlogic-no/packages/linux.5.15-bak/udev.d/30-disable-wakeup.rules
# SPDX-License-Identifier: GPL-2.0-or-later
# Copyright (C) 2009-2014 Stephan Raue (stephan@openelec.tv)
# NextOS 2026-05-12: ${devpath} (aspas simples) era invalido para udev moderno;
# trocado por $env{DEVPATH} (forma canonica). Sem essa correcao a rule inteira
# era ignorada e USB wakeup ficava enabled.

ACTION!="add|change", GOTO="end"

DRIVER=="ehci-pci|xhci_hcd", RUN+="/bin/sh -c 'echo disabled > /sys$env{DEVPATH}/power/wakeup'"

LABEL="end"
