# projects/NextOS/devices/Amlogic-no/packages/linux-drivers/openvfd-driver/system.d/vfd-icons.service
[Unit]
Description=VFD status icons (wifi/eth/usb)
After=openvfd.service
Wants=openvfd.service

[Service]
Type=simple
ExecStart=/usr/lib/coreelec/vfd-icons
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
