# projects/NextOS/packages/network/syncthing/system.d/syncthing.service
[Unit]
Description=Start SyncThing
After=nextos.service
ConditionPathExists=/storage/.cache/services/syncthing.conf

[Service]
Type=simple
ExecStart=/usr/bin/start_syncthing.sh

[Install]
WantedBy=nextos.service
