# projects/NextOS/packages/sysutils/systemd/system.d/systemd-resolved.service.d/override.conf
[Unit]
# systemd-resolved was sitting in the boot critical chain via
# nss-lookup.target → basic.target → sysinit.target, costing ~4.3s in
# blame. None of NextOS's boot-critical services actually need DNS
# (EmulationStation, sway, automount). Drop the implicit pull so
# resolved starts in the background while the UI comes up; DNS is ready
# by the time the user actually opens a network app.
Wants=
Before=systemd-user-sessions.service shutdown.target

[Service]
# Relax sandboxing to prevent race conditions on boot with OverlayFS
# Fixes error: NAMESPACE "Failed to set up mount namespacing" which can occur based on overlayfs setup speed
PrivateTmp=no
ProtectSystem=full
TimeoutStopSec=1s
