# projects/NextOS/packages/nextos/system.d/nextos-report-stats.timer
[Unit]
Description=Timer to report anonymous NextOS install statistics
ConditionPathExists=/storage/.cache/services/stats.enabled

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
