# projects/NextOS/packages/nextos/system.d/nextos-report-stats.service
[Unit]
Description=Report anonymous NextOS install statistics

[Service]
Type=oneshot
ExecStart=/usr/bin/nextos-report-stats

[Install]
WantedBy=multi-user.target
