# projects/NextOS/packages/nextos/system.d/nextos-automount.service
[Unit]
Description=NextOS user autostart script
Before=autostart.service
After=systemd-tmpfiles-setup.service
RequiresMountsFor=/storage

[Service]
Type=oneshot
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/automount
RemainAfterExit=yes

[Install]
WantedBy=nextos.target
