# projects/NextOS/packages/sysutils/systemd/system.d/systemd-hwdb-update.service.d/override.conf
[Unit]
# /usr is read-only squashfs/ext2 in NextOS — the hwdb source files never
# change between boots, so rebuilding hwdb.bin at every boot was costing
# ~4s on RK3326 (systemd-analyze blame: 4.417s) and showed up in the
# critical chain via sysinit.target.
#
# We rebuild hwdb.bin at image creation (see systemd post_makeinstall) and
# ship it with the image. At runtime, skip the rebuild unconditionally by
# ConditionPathExists pointing at a sentinel path that never exists.
ConditionPathExists=/run/systemd/never-create-this

[Service]
# Defensive: if the condition somehow doesn't kick in, exit immediately
# so we still avoid the 4s rebuild.
ExecStart=
ExecStart=/bin/true
