# projects/NextOS/packages/sysutils/busybox/system.d/var-log.mount
[Unit]
Description=Persistent Log Storage
RequiresMountsFor=/var /storage
Requires=storage-log.service
After=storage-log.service storage.mount
ConditionKernelCommandLine=!installer

[Mount]
What=/storage/.cache/log
Where=/var/log
Options=bind
LazyUnmount=yes
